# recipes-core/systemd/obmc-targets/obmc-mpipl@.target
[Unit]
Description= Memory Preserving IPL.
Wants=multi-user.target
After=multi-user.target
OnFailure=obmc-chassis-poweroff@%i.target
