# recipes-phosphor/ipmi/phosphor-ipmi-net/phosphor-ipmi-net@eth1.socket
[Socket]
ListenDatagram=623
BindToDevice=eth1

[Install]
WantedBy=sockets.target

