# recipes-phosphor/ipmi/phosphor-ipmi-net/phosphor-ipmi-net@eth0.socket
[Socket]
ListenDatagram=623
BindToDevice=eth0

[Install]
WantedBy=sockets.target

